# ansible/playground/playbooks/roles/setup_vs_code/files/vscode-server.service
[Unit]
Description=code-server
After=network.target

[Service]
Type=exec
Restart=always
User=%i
ExecStart=/usr/bin/code-server /home/%i/

[Install]
WantedBy=default.target
